# a/UNUSED/aaa_cycojesus/config/etc/profile.d/00_cycojesus.sh
#!/bin/sh

# "Fix" compilation on a multilib system
export LDFLAGS="-L/usr/lib64 "
